# scripts/launchd/team-weekly-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>LABEL_PREFIX.team-weekly-report</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>REPO_PATH/scripts/weekly-report.sh</string>
  </array>
  <!-- Friday (Weekday=5) at 15:00 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>5</integer>
    <key>Hour</key><integer>15</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>REPO_PATH/vault/.logs/launchd-weekly.out</string>
  <key>StandardErrorPath</key>
  <string>REPO_PATH/vault/.logs/launchd-weekly.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
